# assets/midos-house.service
[Unit]
Description=Web server for midos.house
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
ExecStart=/home/fenhl/bin/midos-house
User=mido
Group=nogroup
WorkingDirectory=/opt/git/github.com/midoshouse/midos.house/master

[Install]
WantedBy=multi-user.target
